# datasafe-lite.service
[Service]
Type=oneshot
ExecStart=/usr/sbin/datasafe-lite
